# backend/Database/dailyReset.timer
[Unit]
Description=Runs daily reset sudoku service

[Timer]
OnCalendar=daily
Persistent=True

[Install]
WantedBy=timers.target
