# backend/Database/draftDailyReset.service
[Unit]
Description=Reset daily rankings sudoku table

[Service]
Type=OneShot

