# backend/Database/draftDailyReset.service
[Unit]
Description=Truncate daily rankings sudoku table

[Service]
Type=OneShot

