# backend/Database/draftWeeklyReset.service
[Unit]
Description=Truncate weekly rankings sudoku table
