# backend/Database/draftWeeklyReset.service
[Unit]
Description=Reset weekly rankings sudoku table

[Service]
Type=OneShot

