# backend/Database/weeklyReset.timer
[Unit]
Description=Runs weekly reset sudoku service

[Timer]
OnCalendar=weekly
Persistent=True

[Install]
WantedBy=timers.target

