# backend/Database/weeklyReset.timer
[Unit]
Description=Truncate daily rankings sudoku table

[Timer]
OnCalendar=daily
Persistent=True

[Install]
WantedBy=timers.target

